# bin/served
#!/bin/sh
# author: Daniel (dmilith) Dettlaff

. /etc/rc.subr

SVDSS_BIN="/bin/svdss"
RM_BIN="/bin/rm"
CAT_BIN="/bin/cat"
KILL_BIN="/bin/kill"

name="served"
pid_file="/SystemUsers/.0.pid"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
restart_cmd="${name}_restart"
extra_commands="status"


served_status() {
    if [ -e "${pid_file}" ]; then
        pid="$(${CAT_BIN} ${pid_file})"
        ki="$(${KILL_BIN} -0 ${pid})"
        if [ "$?" = "0" ]; then
            echo "* Service Spawner is running."
            return
        else
            ${RM_BIN} -f "$pid_file"
        fi
    fi
    echo "* Service Spawner is stopped."
}


served_start() {
    if [ -e "${pid_file}" ]; then
        pid="$(${CAT_BIN} ${pid_file})"
        ki="$(${KILL_BIN} -0 ${pid})"
        if [ "$?" = "0" ]; then
            echo "* Service Spawner is already running."
            return
        else
            ${RM_BIN} -f "$pid_file"
        fi
    fi
    echo "* Launching Service Spawner."
    ${SVDSS_BIN} > /SystemUsers/ss.log 2>&1 &
}


served_restart() {
    served_stop
    served_start
}


served_stop() {
    if [ -e "${pid_file}" ]; then
        pid="$(${CAT_BIN} ${pid_file})"
        ki="$(${KILL_BIN} -INT ${pid})"
        if [ "$?" = "0" ]; then
            echo "* Service Spawner has been stopped."
            return
        else
            ${RM_BIN} -f "$pid_file"
        fi
    else
        echo "* Service Spawner is already stopped."
    fi
}


load_rc_config $name
run_rc_command "$1"
