# bin/served
#!/bin/sh

. /etc/rc.subr

export PATH="$PATH:/Software/Openjdk7/bin"
export JAVA_HOME="/Software/Openjdk7/bin"

name="served"
pid_file="/SystemUsers/0.pid"
start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
extra_commands="status"


served_status() {
    if [ -e "$pid_file" ]; then
        root=$(jps | grep -v Jps |grep root\.core | wc -l)
        echo "* RootCore worker:"
        echo "== $(echo ${root})"

        user=$(ps -Af |grep user\.core | grep -v grep | wc -l)
        echo "* UserCore workers:"
        echo "== $(echo $user - 1 | bc)"
    else
        echo "No workers around."
    fi
}


served_start() {
    echo "Starting ServeD KickStart"
    sleep 3
    core=$(ps -Af |grep root\.core | wc -l)
    if [ $core = 2 ]; then
        echo "Already started!"
        return
    else
        svdkick & 2>&1
        return
    fi
}


served_stop() {
    echo "Stopping ServeD KickStart"
    sleep 3
    if [ ! -e "$pid_file" ]; then
        echo "Already stopped!"
        return
    else
        kill -INT "$(cat $pid_file)"
        return
    fi
}


load_rc_config $name
run_rc_command "$1"
