# bin/served
#!/bin/sh

. /etc/rc.subr

export PATH="$PATH:/Software/Openjdk64/bin"
export JAVA_HOME="/Software/Openjdk64"

if [ ! -e $JAVA_HOME/bin/java ]; then
    echo "Cannot find OpenJDK"
    exit 1
fi

name="served"
pid_file="/SystemUsers/0.pid"
start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
restart_cmd="${name}_restart"
extra_commands="status"
root=$(jps | grep -v Jps |grep root\.core | wc -l)
users=$(ps -Af |grep user\.core | grep -v grep | wc -l)


served_status() {
    if [ -e "$pid_file" ]; then
        echo "* RootCore worker:"
        echo "== $(echo ${root})"
        echo "* UserCore workers:"
        echo "== $(echo ${users})"
    else
        echo "No workers around."
    fi
}


served_start() {
    echo "Starting ServeD KickStart"
    core=$(jps | grep -v Jps |grep root\.core | wc -l)
    if [ "${core}" -eq "1" ]; then
        echo "Already started!"
        return
    else
        svdkick > /dev/null 2>&1 &
        echo "ServeD started."
        return
    fi
}


served_restart() {
    served_stop
    for i in $(jot 10); do
        amount=$(ps -Au |grep java | grep -v grep | wc -l)
        if [ "$amount" -gt "0" ]; then
            echo "Some arrogant process found. Sending napalm."
            killall -TERM java
        fi
        if [ "$amount" -eq "0" ]; then
            break
        fi
        sleep 10
    done
    for i in $(jot 10); do
        amount=$(ps -Au |grep java | grep -v grep | wc -l)
        if [ "$amount" -gt "0" ]; then
            echo "Some bastard is still alive! Sending anihilation."
            killall -KILL java
        fi
        if [ "$amount" -eq "0" ]; then
            break
        fi
        sleep 1
    done
    served_start
}


served_stop() {
    echo "Stopping ServeD KickStart"
    if [ "$(jps | grep -v Jps |grep root\.core | wc -l)" -eq "1" ]; then
        kill -INT "$(cat $pid_file)"
        return
    else
        echo "Already stopped!"
        return
    fi
}


load_rc_config $name
run_rc_command "$1"
