# bin/served
#!/bin/sh
# author: Daniel (dmilith) Dettlaff

. /etc/rc.subr
. /etc/profile_sofin

name="served"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
restart_cmd="${name}_restart"
extra_commands="status"


served_start() {
    echo "* Launching SuperUser Service Spawner"
    daemon svdss

    cd /Users
    for user in *; do
        if [ -f "/Users/${user}/.autostart" ]; then
            echo "* Launching Service Spawner for user ${user}"
            uid="$(stat -f%u "/Users/${user}")"
            svdshell -u${uid} -- daemon svdss
        fi
    done

}


served_restart() {
    served_stop
    sleep 1
    served_start
}


served_stop() {
    echo "* Terminating Service Spawner."
    killall svdss
}


load_rc_config $name
run_rc_command "$1"
