# bin/served
#!/bin/sh
# author: Daniel (dmilith) Dettlaff

. /etc/rc.subr
. /etc/profile_sofin

name="served"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
restart_cmd="${name}_restart"
extra_commands="status"


served_start() {
    echo "* Launching Service Spawner."
    nohup svdss &

    # cd /Users
    # for user in *; do
    #     uid="$(find /Users -maxdepth 1 -uid "${CURRENT_USER_UID}" 2> /dev/null)" # get user dir by uid and ignore access errors
    #     echo "* Launching Service Spawner for user: ${user}"
    #     svdshell -u sdfdsfdsfd ${user} -c "svdss >/dev/null 2>&1 &"
    # done
}


served_restart() {
    served_stop
    sleep 1
    served_start
}


served_stop() {
    echo "* Terminating Service Spawner."
    killall svdss
}


load_rc_config $name
run_rc_command "$1"
