# bin/served
#!/bin/sh

. /etc/rc.subr

export PATH="$PATH:/Software/Openjdk7/bin"
export JAVA_HOME="/Software/Openjdk7"

if [ ! -e $JAVA_HOME/bin/java ]; then
    echo "Cannot find OpenJDK"
    exit 1
fi

name="served"
pid_file="/SystemUsers/0.pid"
start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
extra_commands="status"
root=$(jps | grep -v Jps |grep root\.core | wc -l)
users=$(ps -Af |grep user\.core | grep -v grep | wc -l)


served_status() {
    if [ -e "$pid_file" ]; then
        echo "* RootCore worker:"
        echo "== $(echo ${root})"
        echo "* UserCore workers:"
        echo "== $(echo ${users})"
    else
        echo "No workers around."
    fi
}


served_start() {
    echo "Starting ServeD KickStart"
    core=$(jps | grep -v Jps |grep root\.core | wc -l)
    if [ $core = 1 ]; then
        echo "Already started!"
        return
    else
        svdkick > /dev/null 2>&1 &
        echo "ServeD started."
        return
    fi
}


served_stop() {
    echo "Stopping ServeD KickStart"
    sleep 3
    if [ ! -e "$pid_file" ]; then
        echo "Already stopped!"
        return
    else
        kill -INT "$(cat $pid_file)"
        return
    fi
}


load_rc_config $name
run_rc_command "$1"
